# raspberrypi/vscode.service
[Unit]
    Description=VSCode Tunnel Service
    After=network.target

[Service]
    ExecStart=/usr/bin/code tunnel --accept-server-license-terms
    Restart=always

[Install]
    WantedBy=default.target

